# qr3_97.conf
[program:qr3_97]
command=python3 /home/curator3/qr3_97/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_97/test.err.log
stdout_logfile=/home/curator3/qr3_97/test.err.log
